# .quadlet/metron-web.container
[Unit]
Description=Metron Django Web Application
After=metron-postgres.service metron-redis.service
Requires=metron-postgres.service metron-redis.service

[Container]
Image=localhost/metron:latest
ContainerName=metron-web
Network=host
EnvironmentFile=%h/.config/containers/metron.env
LogDriver=journald

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
